# daemon/videonorma.service
[Unit]
Description=videonorma — video loudness normalization daemon
Documentation=https://github.com/prostopasta/videonorma
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/videonorma-daemon
Restart=on-failure
RestartSec=10

# Pass graphical session environment
Environment=DISPLAY=:0
PassEnvironment=DBUS_SESSION_BUS_ADDRESS DISPLAY WAYLAND_DISPLAY XDG_RUNTIME_DIR

[Install]
WantedBy=graphical-session.target
